# bin/mqtt433.service
[Unit]
Description=MQTT to RF433 gateway.
Requires=network.target
After=loxberry.service

[Service]
ExecStart=REPLACELBPBINDIR/mqtt433.sh
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal
#User=loxberry
#Group=loxberry

[Install]
WantedBy=multi-user.target
